# launchd/com.myproject.validator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myproject.validator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/yourname/my-project/scripts/auto-validator.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1500</integer>
  <key>StandardOutPath</key>
  <string>/Users/yourname/my-project/.ai/logs/launchd-validator.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yourname/my-project/.ai/logs/launchd-validator.err.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
